# deploy/systemd/signaldesk-sdrtrunk.service
[Unit]
Description=SignalDesk SDRTrunk
After=signaldesk-fluxbox.service signaldesk-x11vnc.service
Requires=signaldesk-fluxbox.service signaldesk-x11vnc.service

[Service]
Type=simple
User=signaldesk
Group=signaldesk
Environment=DISPLAY=:99
Environment=HOME=/var/lib/signaldesk
ExecStart=/bin/bash /opt/signaldesk/scripts/sdrtrunk-launch.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
